# etc/gentoo/services/getty/run
#!/bin/dash
exec chpst -P /sbin/agetty 38400 TTY linux
