# etc/gentoo/services/lightdm/run
#!/bin/dash
exec 2>&1
echo "starting lightdm..."
echo "waiting for dbus..."
sv start dbus || exit 1
exec lightdm
