# etc/gentoo/services/dbus/run
#!/bin/dash
exec 2>&1
echo "starting dbus..."
echo "waiting for syslog..."
sv start syslog || exit 1
exec dbus-daemon --system --nofork
