# etc/gentoo/services/metalog/run
#!/bin/dash
exec 2>&1
echo "starting metalog..."
exec metalog
