# etc/gentoo/services/ntp/run
#!/bin/dash
exec 2>&1
echo "waiting for syslog..."
sv start syslog || exit 1
echo "starting ntp..."
exec ntpd --nofork
