# etc/gentoo/services/consolekit/run
#!/bin/dash
exec 2>&1
echo "waiting for dbus..."
sv start dbus || exit 1
echo "starting consolekit..."
exec console-kit-daemon --no-daemon
