# etc/gentoo/services/cups/run
#!/bin/dash
exec 2>&1
echo "waiting for syslog..."
sv start syslog || exit 1
echo "starting cups..."
exec cupsd -f
