# etc/gentoo/services/cups/run
#!/bin/dash
. /etc/runit/lib-stage-2.sh
exec 2>&1
wait_for syslog
echo "starting cups..."
exec cupsd -f
