# etc/gentoo/services/udev/run
#!/bin/dash
. /etc/runit/lib-stage-2.sh
wait_for dbus
echo "starting udev..."
exec udev
